# meta-pitv/recipes-multimedia/pitv-bootsplash/files/pitv-bootsplash.service
[Unit]
Description=PiTv boot splash video
# Wait for Plymouth to fully release the DRM device before we take it.
After=plymouth-quit-wait.service
Wants=plymouth-quit-wait.service
Before=display-manager.service sddm.service

[Service]
Type=oneshot
ExecStart=/usr/bin/pitv-bootsplash
TimeoutStartSec=120
# Acquire VT1 as foreground so mpv can open the DRM device.
StandardInput=tty
StandardOutput=journal
TTYPath=/dev/tty1

[Install]
WantedBy=graphical.target
